# dist/supervisor.conf
[supervisord]
nodaemon=true

[program:dependencywatcher-tasks]
startretries=10
user=www-data
command=celery -A dependencywatcher.tasks.celery worker --beat --workdir=/var/www/dependencywatcher -s /var/lib/dw-workspace/celery-schedule

[program:dependencywatcher-web]
startretries=10
user=www-data
directory=/var/www/dependencywatcher
command=/var/www/dependencywatcher/run.sh

